# bin/pipower3.service
# https://www.freedesktop.org/software/systemd/man/systemd.service.html
[Unit]
Description=pipower3 service
# Need to start last to avoid gpio being occupied
After=multi-user.target

[Service]
Type=forking
# WorkingDirectory=/opt/pipower3
ExecStart=/usr/local/bin/pipower3 start --background
ExecStop=/usr/local/bin/pipower3 stop
# PrivateTmp=False

[Install]
WantedBy=multi-user.target
